# nova-cert.conf
[program:nova-api]
directory=/
command=/usr/bin/nova-cert --config-file=/etc/nova/nova.conf --log-file=/var/log/nova/nova-cert.log
user=nova
autostart=true
autorestart=true
priority=999
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log